# 50-rtc-hwclock.rules
# /etc/udev/rules.d/50-rtc-hwclock.rules

KERNEL=="rtc0", RUN+="/sbin/hwclock -s"
