# rtc-i2c.rules
# /etc/udev/rules.d/rtc-i2c.rules

KERNEL=="rtc0", RUN+="/sbin/hwclock --hctosys"
